# ops/systemd/muslim-backup-verify.service
[Unit]
Description=Verify MuslimKit encrypted backups
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=muslimkit
Group=muslimkit
WorkingDirectory=/opt/muslimkit/api
EnvironmentFile=/etc/muslimkit/muslim-users-api.env
ExecStart=/usr/bin/bun scripts/verify-backups.ts --dir /var/backups/muslimkit/api --limit 24
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
MemoryDenyWriteExecute=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=027
ReadWritePaths=/opt/muslimkit/api /var/backups/muslimkit/api /var/log/muslimkit
